# monitor-agent/systemd/frogstats-monitor-agent.service
[Unit]
Description=FrogStats Monitor Agent
Documentation=https://github.com/frogstats/monitor-agent
After=network.target

[Service]
Type=simple
User=frogstats
Group=frogstats
WorkingDirectory=/opt/frogstats-monitor-agent
ExecStart=/usr/bin/node /opt/frogstats-monitor-agent/dist/index.js
Restart=always
RestartSec=10

# Environment
Environment="NODE_ENV=production"
EnvironmentFile=/opt/frogstats-monitor-agent/.env

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/frogstats-monitor-agent/logs

# Resource limits
LimitNOFILE=65536
MemoryLimit=1G
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=frogstats-monitor-agent

[Install]
WantedBy=multi-user.target